# HOME/.crontab
# +----------------------- minute (0 - 59)
# |    +------------------ hour (0 - 23)
# |    |    +------------- day of month (1 - 31)
# |    |    |    +-------- month (1 - 12)
# |    |    |    |    +--- day of week (0 - 7) (Sunday=0 or 7)
# |    |    |    |    |
# *    *    *    *    *      command to be executed
 */5   *    *    *    *      if [ -e $HOME/.procimap/notify.py ]; then $HOME/.procimap/notify.py > /dev/null; fi
 38    *    *    *    *      if [ -e $HOME/bin/actlog_check_status.pl ]; then $HOME/bin/actlog_check_status.pl; fi
 0    2,13  *    *    *      $HOME/.pyenv/versions/3.8.10/bin/python $HOME/.dotfiles/deploy.py --quiet
